# buildroot-2016.02/board/9tripod/x4418/rootfs-overlay/etc/profile.d/qt5.sh
export QT_QPA_PLATFORM_PLUGIN_PATH=/usr/lib/qt/plugins/platforms
export QT_QPA_PLATFORM=linuxfb:fb=/dev/fb0
export QT_QPA_FONTDIR=/usr/share/fonts
export QT_QPA_GENERIC_PLUGINS=evdevkeyboard,evdevmouse,evdevtouch
